# diggy-repl/docker/socket.conf
[program:socket]
command=/usr/bin/node /opt/app/server.js
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true