# diggy-repl/docker/web.conf
[program:web]
command=/usr/bin/node /opt/app/build/index.js
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true