# diggy-sandbox/docker/web.conf
[program:web]
command=/root/.local/bin/gunicorn -c /sandbox/config/gunicorn.conf.py sandbox.api:app
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true